# debian/kernel-autoremove.upstart
# kernel autoremove - aurex kernel cleanup job
#
# kernel-autoremove removes unused obsolete kernel and headers
# packages to free disk space (specially in /boot)

description     "kernel autoremove"
author          "Luis Antonio Garcia Gisbert<luisgg@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]

task

exec /usr/sbin/kernel-autoremove remove
