# debian/kernel-autoremove.cron.daily
#!/bin/sh
/usr/sbin/kernel-autoremove remove || true
